# etc/hotplug.d/ntp/post-ntpd-event
#!/bin/sh

/usr/lib/node_modules/@yoda/flora/bin/flora-cli post yodaos.ntpd.event "[\"$ACTION\"]"
